# cookbooks/clsi/templates/default/cron.d/check_clsi.cron
*/5 * * * * root ruby <%= node[:clsi][:install_directory] %>/shared/check.rb
