# clsi/templates/default/clean_output_and_cache.cron
#!/bin/sh
cd <%= node[:clsi][:install_directory] %>/current
rake clsi:clean_cache RAILS_ENV=production
rake clsi:clean_output RAILS_ENV=production

